# 99-sleep-on-unplug.rules
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/home/seth/git/wake-on-plug/sleep-on-unplug.sh"
