# 99-wake-on-plug.rules
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/home/seth/git/wake-on-plug/wake-on-plug.sh"
